# deployment/waverover.service.d/10-auto-update.conf
[Unit]
Wants=waverover-update.service
After=waverover-update.service

[Service]
Environment=WAVEROVER_IDENTITY_FILE=/home/waverover/ros2_ws/src/waverover/config/robot_identity.yaml
